# scripts/launchd/com.constructpc.miniapp.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.constructpc.miniapp.web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexander/Documents/Construct/.nodeenv/bin/node</string>
    <string>/Users/alexander/Documents/Construct/miniapp_web/server.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexander/Documents/Construct/miniapp_web</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOST</key>
    <string>0.0.0.0</string>
    <key>PORT</key>
    <string>8081</string>
    <key>CONSTRUCT_API_ORIGIN</key>
    <string>http://127.0.0.1:8080</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alexander/Documents/Construct/data/launchd.miniapp_web.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexander/Documents/Construct/data/launchd.miniapp_web.err.log</string>
</dict>
</plist>
