# scripts/launchd/com.constructpc.miniapp.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.constructpc.miniapp.api</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexander/Documents/Construct/.venv/bin/python</string>
    <string>-m</string>
    <string>miniapp_api.run_service</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexander/Documents/Construct</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>MINIAPP_APPLY_MIGRATIONS</key>
    <string>1</string>
    <key>MINIAPP_RELOAD</key>
    <string>0</string>
    <key>MINIAPP_API_HOST</key>
    <string>0.0.0.0</string>
    <key>MINIAPP_API_PORT</key>
    <string>8080</string>
    <key>PYTHONPYCACHEPREFIX</key>
    <string>/Users/alexander/Documents/Construct/.pycache_tmp</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alexander/Documents/Construct/data/launchd.miniapp_api.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexander/Documents/Construct/data/launchd.miniapp_api.err.log</string>
</dict>
</plist>
