# build/ferrum_release/scripts/launchd/com.constructpc.miniapp.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.constructpc.miniapp.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexander/Documents/Construct/.venv/bin/python</string>
    <string>/Users/alexander/Documents/Construct/scripts/run_miniapp_tunnel.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexander/Documents/Construct</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PYTHONPYCACHEPREFIX</key>
    <string>/Users/alexander/Documents/Construct/.pycache_tmp</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alexander/Documents/Construct/data/launchd.miniapp_tunnel.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexander/Documents/Construct/data/launchd.miniapp_tunnel.err.log</string>
</dict>
</plist>
